# roles/todo/templates/todo.service
[Unit]
Description = todo Service
[Service]

Environment=ZIPKIN_URL=
Environment=REDIS_HOST=http://redis.zssurendra.online:8080
Environment=TODOS_API_PORT=8080
ExecStart=/bin/node /root/server.js
SyslogIdentifier=login

[Install]
WantedBy=multi-user.target